# crontab
0 * * * * root . /etc/environment && php /app/refurbed-invoices.php >> /proc/1/fd/1 2>&1
